# templates/default/tomcat/cleaner.cron.erb
# Clean cache files not used by 30 minutes
*/<%= node['tomcat']['cleaner.minutes.interval'] %> * * * * root find <%= @tomcat_cache_path %> -mmin +<%= node['tomcat']['cleaner.minutes.interval'] %> -type f -exec rm {} \;

# Clean rotated logs
20 0 * * * root rm -rf <%= @tomcat_log_path %>/alfresco.*.log
# TODO - this also removes the currently used one; fix it
# 20 0 * * * root rm -rf <%= @tomcat_log_path %>/catalina.*.log
# 20 0 * * * root rm -rf <%= @tomcat_log_path %>/localhost.*.log
