# templates/default/tomcat/lsof_tomcat.cron.erb
MAILTO=""

0 * * * * root /usr/sbin/lsof -u tomcat >> <%= @tomcat_log_path %>/lsof_tomcat.log
